# crontab
0 7-13 * * 1-5 ./spit/submit.mjs outfield $(./outfield/index.mjs -e)
