# crontab
31 6 * * 1-5 root ./spit/submit.mjs outfield $(./outfield/index.mjs)
