# crontab
31 6 * * 1-5 ./spit/submit.mjs outfield $(./outfield/index.mjs)
